# bin/ccolejorz.service
[Unit]
Description=ccolejorz stationmaster

[Service]
Type=simple
ExecStart=/home/pi/.virtualenvs/colejorz/bin/colejorz
RemainAfterExit=no
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=ccolejorz

[Install]
WantedBy=multi-user.target
